# mock_systemd_units/system-monitor.service
[Unit]
Description=System Monitoring Agent

[Service]
ExecStart=/usr/bin/monitor-agent
Type=simple

[Install]
WantedBy=multi-user.target
